# photomanager.conf
[program:photomanager]
command: /usr/bin/python manage.py run_gunicorn
directory: /home/ubuntu/
autostart: true
environment=DJANGO_SETTINGS_MODULE="/home/ubuntu/PhotoApp/dev_settings.py"