# photomanager.conf
[program:photomanager]
command: /usr/bin/python -m microblog
directory: /home/ubuntu/FlaskMicroblog
autostart: true
environment=DJANGO_SETTINGS_MODULE="/home/ubuntu/PhotoApp/dev_settings.py"