# photomanager.conf
[program:photomanager]
command: gunicorn PhotoApp.wsgi:application
directory: /home/ubuntu/DjangoApp
autostart: true
environment=DJANGO_SETTINGS_MODULE="PhotoApp.dev_settings"